# startup/com.tangynature.quote.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tangynature.quote</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/harsimransingh/Documents/GitHub/tangy-nature-quote/Tangy_Nature_Quote.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/harsimransingh/Library/Logs/tangynature-quote.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/harsimransingh/Library/Logs/tangynature-quote-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/harsimransingh/Documents/GitHub/tangy-nature-quote</string>
</dict>
</plist>
